# components/example-notebook-servers/jupyter/s6/services.d/sshd/run
#!/command/with-contenv bash

exec sudo /usr/sbin/sshd -D